# roles/hdmi-media/systemd/mpv-hdmi@.service
[Unit]
Description=MPV HDMI media playback (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=video
Environment=HDMI_CONNECTOR=HDMI-A-1
Environment=HDMI_AUDIO_DEVICE=plughw:vc4hdmi,0
ExecStart=/usr/bin/mpv --idle=yes --fs --gpu-context=drm --drm-connector=${HDMI_CONNECTOR} --audio-device=alsa/${HDMI_AUDIO_DEVICE} --input-ipc-server=/run/mpv.sock --no-terminal
Restart=always

[Install]
WantedBy=multi-user.target

