# agent/role-agent.timer
[Unit]
Description=Fleet GitOps converger timer
Requires=role-agent.service

[Timer]
OnBootSec=45s
OnUnitActiveSec=10m
RandomizedDelaySec=90
Persistent=true

[Install]
WantedBy=timers.target
