# agent/role-agent-watchdog.service
[Unit]
Description=Fleet agent convergence watchdog
After=role-agent.service

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /opt/fleet/agent/role-agent-watchdog.sh

[Install]
WantedBy=multi-user.target
