# agent/role-agent-watchdog.timer
[Unit]
Description=Fleet agent convergence watchdog timer
Requires=role-agent-watchdog.service

[Timer]
OnBootSec=3m
OnUnitActiveSec=5m
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
